# poky/meta-HU/recipes-client/src/files/client.service
[Unit]
Description=My Application
After=network.target

[Service]
ExecStart=/usr/bin/HelloWorldClient
WorkingDirectory=/usr/bin
Restart=always

[Install]
WantedBy=multi-user.target

